# gentoo-init-script
#!/sbin/runscript

depend() {
        need net
}

start() {
        ebegin "Starting tornado"
        start-stop-daemon --start -i --user apache:apache --quiet --background --make-pidfile -                                 
-pidfile /var/run/tornado.pid  --exec /usr/bin/python /web/crypto.is/rss.py
        eend $?
}

stop() {
        ebegin "Stopping tornado"
        start-stop-daemon --stop --pidfile /var/run/tornado.pid
        eend $?
}

